# package/src/deb-scripts/postinst
cp /usr/local/share/vdriver/vdriver.service /lib/systemd/system/vdriver.service
chmod 644 /lib/systemd/system/vdriver.service
systemctl daemon-reload
systemctl enable vdriver.service
systemctl start vdriver.service
systemctl status vdriver.service --no-pager || true