# package/src/deb-scripts/prerm
systemctl disable vdriver.service 2>/dev/null
systemctl stop vdriver.service 2>/dev/null